# rluks-auto-cleanup.service
[Unit]
Description=Automatic cleanup

[Service]
Type=oneshot
ExecStart=/home/roman/Documents/GitHub/auto-cleanup/cleanup.sh
