# rluks-auto-cleanup.timer
[Unit]
Description=Daily automatic file cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
